# contrib/lite-google-workspace-mcp@.service
[Unit]
Description=lite-google-workspace-mcp (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/Codes/lite-google-workspace-mcp/.venv/bin/lite-google-workspace-mcp serve --account %i
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
